# dotfiles/launchd/de.gershu.nova.lab.ib.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  nova-lab IB Gateway LaunchAgent — laeuft in stefan_mac's GUI-Session.

  Gateway ist eine GUI-Java-App und braucht WindowServer-Zugriff. Das geht
  NICHT mit einem system-Daemon (egal welcher UserName-Key gesetzt ist —
  Aqua haengt am gui/UID-Domain, nicht am Prozess-Owner). Deshalb:
  LaunchAgent unter stefan_mac, der auf nova-hub permanent eingeloggt ist.

  Initial-Setup (als stefan_mac, mit aktiver GUI-Session):
    cp ~/Path/to/nova/dotfiles/launchd/de.gershu.nova.lab.ib.gateway.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) \
       ~/Library/LaunchAgents/de.gershu.nova.lab.ib.gateway.plist
    launchctl kickstart -k gui/$(id -u)/de.gershu.nova.lab.ib.gateway

  Logs:    /Users/stefan_mac/Library/Logs/nova-ib-gateway.log
  Status:  launchctl print gui/$(id -u)/de.gershu.nova.lab.ib.gateway
  Stop:    launchctl bootout gui/$(id -u)/de.gershu.nova.lab.ib.gateway
  Restart: launchctl kickstart -k gui/$(id -u)/de.gershu.nova.lab.ib.gateway
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.gershu.nova.lab.ib.gateway</string>

    <!-- KEIN UserName/GroupName-Key: LaunchAgents laufen automatisch als
         der User, in dessen gui/UID-Domain sie geladen werden. -->


    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>exec /Users/novaadm/nova/scripts/ib_gateway_start.sh</string>
    </array>

    <!-- ON-DEMAND: kein RunAtLoad/KeepAlive. Das Gateway belegt den IBKR-
         Account NICHT permanent (damit direkte TWS-/Mobile-Nutzung moeglich
         bleibt). Start nur bei Bedarf:
           scripts/ib_gateway_ctl.sh start   (intern launchctl kickstart)
         Stop:  scripts/ib_gateway_ctl.sh stop
         Der Agent bleibt geladen, laeuft aber nur nach explizitem kickstart
         und bleibt nach 'stop' aus, bis er erneut angefordert wird. -->
    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <!-- Mindestens 60s zwischen Restart-Versuchen — verhindert
         API-Hammering bei kaputter Config oder waehrend IBKR-Wartung. -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/stefan_mac/Library/Logs/nova-ib-gateway.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/stefan_mac/Library/Logs/nova-ib-gateway.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/stefan_mac</string>
    </dict>

    <!-- WorkingDirectory auf stefan_mac's Home — IBC's mktemp/Session-ID
         braucht ein beschreibbares CWD. /Users/novaadm/nova ginge nicht,
         da stefan_mac dort nur Lese-/Exec-Rechte hat. -->
    <key>WorkingDirectory</key>
    <string>/Users/stefan_mac</string>
</dict>
</plist>
