# dotfiles/launchd/de.gershu.nova.lab.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  nova-lab daily-digest LaunchDaemon fuer nova-hub.

  Triggert taeglich um 23:30 lokal (15 Min nach lab.monitor) als novaadm.
  Setup via scripts/install_lab_digest_daemon.sh (sudo, einmalig).

  Logs: /Users/novaadm/Library/Logs/nova-lab-digest.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.gershu.nova.lab.digest</string>

    <key>UserName</key>
    <string>novaadm</string>

    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/novaadm/nova/scripts/lab_digest_daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-digest.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-digest.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/novaadm</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/novaadm</string>
</dict>
</plist>
