# dotfiles/launchd/de.gershu.nova.lab.superinvestors.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  nova-lab superinvestors LaunchDaemon — 13F-Filer-Ingest.

  Schedule: woechentlich, Sonntag 07:00 (StartCalendarInterval). 13F-Filings
  erscheinen ~45 Tage nach Quartalsende; ein woechentlicher Lauf faengt neue
  Filings ein, sobald sie da sind. Upsert ist idempotent (re-run aktualisiert
  nur). sec-api-lastig, daher selten und getrennt vom 5-Minuten-Worker.

  Initial-Setup: scripts/install_daemon.sh lab.superinvestors
  Logs:          /Users/novaadm/Library/Logs/nova-lab-superinvestors.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.gershu.nova.lab.superinvestors</string>

    <key>UserName</key>
    <string>novaadm</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/novaadm/nova/scripts/lab_superinvestors.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>7</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-superinvestors.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-superinvestors.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/novaadm</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/novaadm/nova</string>
</dict>
</plist>
