# dotfiles/launchd/de.gershu.nova.lab.portfolio_briefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  nova-lab daily portfolio_briefing LaunchDaemon fuer nova-hub.

  Sequenz: monitor 23:15 -> alert_explainer 23:20 -> THIS 23:25 -> digest 23:30.

  Synthetisiert per-currency totals, top-movers, alerts (mit Erklaerungen)
  -> LLM-call auf nova-w5 -> sig_portfolio_briefings.

  Initial-Setup: scripts/install_lab_portfolio_briefing_daemon.sh.

  Logs: /Users/novaadm/Library/Logs/nova-lab-portfolio-briefing.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.gershu.nova.lab.portfolio_briefing</string>

    <key>UserName</key>
    <string>novaadm</string>

    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/novaadm/nova/scripts/lab_portfolio_briefing_daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>25</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-portfolio-briefing.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-portfolio-briefing.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/novaadm</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/novaadm</string>
</dict>
</plist>
