# dotfiles/launchd/de.gershu.nova.lab.alert_explainer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  nova-lab daily alert_explainer LaunchDaemon fuer nova-hub.

  Triggert taeglich um 23:20 Lokalzeit — zwischen monitor (23:15) und
  digest (23:30). Liest sig_alerts vom Tag, holt yfinance-News, ruft LLM
  auf nova-w5 zur Erklaerung, schreibt sig_alert_explanations.

  Initial-Setup: scripts/install_lab_alert_explainer_daemon.sh.

  Logs: /Users/novaadm/Library/Logs/nova-lab-alert-explainer.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.gershu.nova.lab.alert_explainer</string>

    <key>UserName</key>
    <string>novaadm</string>

    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/novaadm/nova/scripts/lab_alert_explainer_daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>20</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-alert-explainer.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-alert-explainer.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/novaadm</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/novaadm</string>
</dict>
</plist>
