# dotfiles/launchd/de.gershu.nova.picker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  nova-picker LaunchDaemon fuer nova-hub (MacMini, headless).

  Wird als LaunchDaemon installiert (NICHT LaunchAgent), weil der Hub
  typischerweise via SSH bedient wird und keine permanente GUI-Login-Session
  hat. LaunchAgents wuerden in dem Setup nie laden — die user/<uid>-Domain
  existiert nur bei aktivem User-Login.

  LaunchDaemon laeuft im System-Kontext, immer aktiv (auch beim Boot,
  ohne dass jemand eingeloggt sein muss). Der UserName-Key sorgt dafuer
  dass der Process trotzdem als novaadm laeuft, nicht als root.

  Initial-Setup: scripts/install_picker.sh (einmalig, mit sudo).
  Liegt in /Library/LaunchDaemons/ (system-weit, root-owned).

  Polling-Intervall: 10 Sekunden. Bei Bedarf in StartInterval anpassen.
  Logs: /Users/novaadm/Library/Logs/nova-picker.log.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.gershu.nova.picker</string>

    <key>UserName</key>
    <string>novaadm</string>

    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/novaadm/nova/scripts/nova_picker.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>10</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/novaadm/Library/Logs/nova-picker.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/novaadm/Library/Logs/nova-picker.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/novaadm</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/novaadm</string>
</dict>
</plist>
