# dotfiles/launchd/de.gershu.nova.lab.screener_csp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  nova-lab daily CSP-Screener LaunchDaemon fuer nova-hub.

  Schedule: 23:05 (nach ingest 23:00, vor monitor 23:15). Nutzt today's
  frisch-ingested spot-quotes + IB option-chain. Ueberschreibt
  'system_recommendations' watchlist mit Top-N CSP-Kandidaten — damit
  digest 23:30 die frischen Picks anzeigen kann.

  Initial-Setup: scripts/install_lab_screener_csp_daemon.sh.

  Logs: /Users/novaadm/Library/Logs/nova-lab-screener-csp.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.gershu.nova.lab.screener_csp</string>

    <key>UserName</key>
    <string>novaadm</string>

    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/novaadm/nova/scripts/lab_screener_csp_daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-screener-csp.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/novaadm/Library/Logs/nova-lab-screener-csp.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/novaadm</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/novaadm</string>
</dict>
</plist>
